# supervisord.conf
[supervisord]
nodaemon=true

[program:caddy]
command=/usr/bin/caddy --conf /etc/Caddyfile --log stdout

[program:go]
command=/bin/api -config /etc/production.yaml